# supervisord.ini
[supervisord]
nodaemon=true

[program:jenkins]
command=/sbin/tini -- /usr/local/bin/jenkins.sh

[program:dockerd]
command=dockerd-entrypoint.sh
